# systemd/plm-sync.service
[Unit]
Description=Sync Designflow PLM master data into shared Supabase
Wants=network-online.target
After=network-online.target
# Loud alert when the sync exits non-zero (see systemd/plm-sync-alert.service).
OnFailure=plm-sync-alert.service

[Service]
Type=oneshot
User=ai
WorkingDirectory=/worksp/shared-db
EnvironmentFile=/home/ai/.plm-sync.env
ExecStart=/usr/bin/bash /worksp/shared-db/tools/run-plm-master-data-sync.sh
